# systemd/fo-secmon.timer
[Unit]
Description=Foundation Security Monitor — every 15 minutes

[Timer]
OnBootSec=4min
OnUnitActiveSec=15min
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target
