# systemd/fo-vpsoc.timer
[Unit]
Description=Foundation VPSOC Dark Sentinel — every 15 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=15min
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target
